# scripts/andresanz-deploy.service
[Unit]
Description=andresanz.com webhook deploy listener
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/technologyfoc.us
ExecStart=/usr/bin/node scripts/webhook-deploy.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/var/www/technologyfoc.us/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
